# dystrybucja/livecd/devel/slax/rootcopy/etc/rc.d/rc.local
#!/bin/sh
/usr/sbin/dhcpd
/usr/share/apache-tomcat/bin/startup.sh

